# root/etc/supervisor.d/elasticsearch.ini
[program:elasticsearch]
command=sudo -u elasticsearch /opt/elasticsearch/bin/elasticsearch
