# root/etc/supervisor.d/php-fpm.ini
[program:php-fpm]
command=/usr/bin/php-fpm --nodaemonize --fpm-config /etc/php/php-fpm.conf
